# bin/private/docker_host_comm/docker_host_comm.service
[Unit]
Description=Docker Host Communication
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/opt/docker_webapps/bin/private/docker_host_comm/docker_host_comm.sh

[Install]
WantedBy=multi-user.target
