# bin/private/docker_host_comm/docker_host_comm.service
[Unit]
Description=Docker Host Communication
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=_BASE_PATH_/bin/private/docker_host_comm/docker_host_comm _BASE_PATH_ 

[Install]
WantedBy=multi-user.target
